# gin_project.conf
[program:gin_project]
command=/root/Go/src/gin_project/main
directory=/root/Go/src/gin_project/
autostart=true
autorestart=true